# supervisor-ini/redispapa.ini
[program:redispapa]
command=/root/www/redispapa/redispapavenv/bin/gunicorn --worker-class socketio.sgunicorn.GeventSocketIOWorker run:app -b 0.0.0.0:8090
directory=/root/www/redispapa
startsecs=0
stopwaitsecs=0
autostart=true
autorestart=true
stdout_logfile=/tmp/redispapa.log
stderr_logfile=/tmp/redispapa.err
