# supervisor-ini/webhooks.ini
[program:webhooks]
command=python index.py 8080
directory=/root/www/github-webhook-handler/
startsecs=0
stopwaitsecs=0
autostart=true
autorestart=true
stdout_logfile=/tmp/webhooks.log
stderr_logfile=/tmp/webhooks.err
