# supervisor-ini/xueer.ini
[program:xueer]
command=gunicorn -w 4 -b 121.41.6.148:5050 wsgi:app
directory=/root/www/xueer_be
startsecs=0
stopwaitsecs=0
autostart=true
autorestart=true
stdout_logfile=/tmp/gunicorn_xueer.log
stderr_logfile=/tmp/gunicorn_xueer.err
